# fs/skeleton-icdtcp3/etc/profile.d/icd.sh
export ICD_DATA_DB=/home/data/data.db
export ICD_CONFIG_DB=/home/data/config.db
export ICD_LIVE_DB=/tmp/live.db

export ICD_TEMP_DIR=/tmp
export ICD_SHARE_DIR=/usr/share/icd
export ICD_DATA_DIR=/home/data



